# systemd/azblogfilter-notify.service
[Unit]
Description=Notify when a new Azure blog post is released

[Service]
Type=simple
